# lcdctrl/95-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", RUN+="/bin/chgrp pi %S%p/brightness", RUN+="/bin/chmod g+w %S%p/brightness"
